# install/forti-hole.service
[Unit]
Description=forti-hole DNS blocklist service
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
OnBootSec=60
ExecStart= # This line is auto-generated via the install script

[Install]
WantedBy=multi-user.target
